# mlhkrc/shell/crontab
0 5 * * * /home/ubuntu/bootcamp/bootcamp/mlhkrc/shell/guncelle.sh >/dev/null 2>&1
