# proje/linux/server-update.service
[Unit]
Description=Update System Service
After=network.target

[Service]
ExecStart=/usr/local/bin/update_system.sh
User=root
Restart=on-failure

[Install]
WantedBy=multi-user.target